# prog2_cs695/rcinit/rc.local
#!/bin/sh -e
#
# rc.local
#
# gcc /home/prateek/server.c -lpthread -o /home/prateek/server &
#./server &
home/prateek/server/server &

exit 0
